# ansible_maintenance_scripts/templates/hdfs-datanode.service.j2
[Unit]
Description=HDFS Namenode Service
After=network.target

[Service]
Type=forking
User="{{ gridgain_rw_user }}"
Environment=HADOOP_HOME={{ hadoop_home }}
Environment=HADOOP_CONF_DIR={{ hadoop_conf_dir }}
Environment=HADOOP_LOG_DIR={{ hadoop_log_dir }}
ExecStart={{ hadoop_home }}/bin/hdfs --config $HADOOP_CONF_DIR --daemon start datanode
ExecStop={{ hadoop_home }}/bin/hdfs --config $HADOOP_CONF_DIR --daemon stop datanode
Restart=on-failure

[Install]
WantedBy=multi-user.target
